# optee/3.10.0_rpi3_netboot_nfs/dhclient.service
[Unit]
Description=enable dhclient

[Service]
Type=simple
ExecStart=/sbin/dhclient eth0
Restart=on-failure

[Install]
WantedBy=multi-user.target
